# deploy/ops/personal-learning-backup.timer
[Unit]
Description=Daily personal learning site backup

[Timer]
OnCalendar=*-*-* 03:30:00
RandomizedDelaySec=10min
Persistent=true

[Install]
WantedBy=timers.target
